# slack-relay/launchd/slack-relay.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: install.sh renders this into ~/Library/LaunchAgents with the
     current user's HOME/venv/repo paths. Nothing user-specific is committed. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__SCRIPT__</string>
        <string>--config</string>
        <string>__CONFIG__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- One-shot poller, not a persistent server: launchd re-fires this
         every StartInterval seconds rather than keeping one process alive
         (mirrors com.sjaconette.cedarmillmonitor's shape, not
         com.agent-console's KeepAlive shape). -->
    <key>StartInterval</key>
    <integer>__INTERVAL__</integer>
    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
