# agent-console/launchd/agent-console.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: install.sh renders this into ~/Library/LaunchAgents with the
     current user's HOME/python/paths. Nothing user-specific is committed. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__SCRIPT__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>SKILLS_DASHBOARD_PORT</key>
        <string>__PORT__</string>
        <key>SKILLS_DASHBOARD_HOST</key>
        <string>__HOST__</string>
        <!-- Dispatch runtime (workboard actions): the `claude` binary a
             dispatch launches is resolved at dispatch time as
             AGENT_CONSOLE_CLAUDE_BIN, else `claude` on PATH, else
             ~/.local/bin/claude. launchd's PATH usually omits ~/.local/bin, so
             set AGENT_CONSOLE_CLAUDE_BIN here if `claude` is not on __PATH__:
               <key>AGENT_CONSOLE_CLAUDE_BIN</key><string>/path/to/claude</string>
             Dispatch logs + records default to
             ~/Library/Logs/agent-console/dispatch (override:
             AGENT_CONSOLE_DISPATCH_DIR). -->
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
