# agentprof/launchd/com.sjaconette.agentprof-refresh.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: the attended install task renders this into
     ~/Library/LaunchAgents with the current user's HOME/script paths.
     Nothing user-specific is committed. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sjaconette.agentprof-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SCRIPT__</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/refresh.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/refresh.log</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
